# rustc-internal-docs.service
[Unit]
Description = Generate Rustc internal docs

[Service]
ExecStart=/home/michael/.cargo/bin/rustc-internal-docs -v
User=michael
